# systemd/jetstream-pg-writer.service
[Unit]
Description=JetStream Postgres Writer
After=docker.service network-online.target
Requires=docker.service
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
Restart=on-failure
RestartSec=10
WorkingDirectory=/home/user/jetstream-pg-writer
ExecStart=/usr/bin/docker compose up --build --detach
ExecStop=/usr/bin/docker compose down
TimeoutStartSec=300
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
